# deploy/synaplex-conjecture-eval-capture.service
[Unit]
Description=Capture production Synaplex conjecture traces for prompt evaluation

[Service]
Type=oneshot
WorkingDirectory=/opt/workspace/projects/synaplex
ExecStart=/bin/bash /opt/workspace/projects/synaplex/scripts/capture-conjecture-eval.sh
